# smtptoxmpp.service
[Service]
ExecStart=/usr/local/sbin/smtptoxmpp /etc/smtptoxmpp.conf